# debian/rootfs/etc/s6-overlay/s6-rc.d/syslog/run
#!/bin/sh

set -e

exec /sbin/syslogd -n
